# ctp_access/quotation/access/ctp/supervisor.conf
[program:ctp_access]
directory = /usr/local/soldier/apps/ctp_access
command = /usr/local/soldier/apps/ctp_access/ctp_access
autostart = true
startsecs = 5 
autorestart = true
startretries = 3
user = epower
redirect_stderr = true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 20
stdout_logfile = /data/applog/apps/ctp_access/stdout.log

