# race_store/java/supervisor.conf
[program:quotation_race_store]
directory = /usr/local/soldier/apps/quotation_race_store/
command = java -server -XX:+PrintGCDetails -XX:+PrintGCDateStamps -jar /usr/local/soldier/apps/quotation_race_store/quotation_race_store-0.0.1-SNAPSHOT-jar-with-dependencies.jar
autostart = true
startsecs = 5 
autorestart = true
startretries = 3
user = epower
redirect_stderr = true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 20
stdout_logfile = /data/applog/apps/quotation_race_store/stdout.log

